# autostart/facebook.desktop
[Desktop Entry]
Name=Facebook
Exec=while :; do facebook-links; sleep 2h; done
Type=Application
X-GNOME-Autostart-enabled=true
