# autostart/facebook.desktop
[Desktop Entry]
Name=Facebook
Exec=facebook-links
Type=Application
X-GNOME-Autostart-enabled=true
