# autostart/give-me-5.desktop
[Desktop Entry]
Name=give-me-5-kurulus-usman
Exec=firefox https://www.giveme5.co/01-20kur
Type=Application
X-GNOME-Autostart-enabled=true
