# ubuntu-18.04-new-os-installations/autostart/lfd mail.desktop
[Desktop Entry]
Name=LFD mail
Exec=https://mail.google.com/mail/u/1/
Type=Application
X-GNOME-Autostart-enabled=true
