# ubuntu-18.04-new-os-installations/autostart/outlook.desktop
[Desktop Entry]
Name=Outlook
Exec=firefox https://outlook.live.com/mail/0/
Type=Application
X-GNOME-Autostart-enabled=true
