# ubuntu-18.04-new-os-installations/autostart/stack_overflow.desktop
[Desktop Entry]
Name=Stack overflow
Exec=firefox https://stackoverflow.com/questions/tagged/pandas+python
Type=Application
X-GNOME-Autostart-enabled=true
